# scripts/hyperagent/com.hyperagent.evolve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hyperagent.evolve</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/choegihwan/Documents/Projects/claude-setup/scripts/hyperagent/evolve.py</string>
    <string>--json</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/choegihwan/Documents/Projects/claude-setup</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/choegihwan/.claude/hyperagent/evolve.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/choegihwan/.claude/hyperagent/evolve.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
